# files/home/.config/systemd/user/night-mode.service
[Unit]
Description=Night Mode Daemon

[Service]
Type=simple
ExecStart=%h/.local/bin/night-mode
Restart=always

[Install]
WantedBy=default.target
